# .rpm/scripts/docuum.service
[Unit]
Description="Docuum LRU Docker image cache manager"

[Service]
ExecStart=/usr/bin/docuum

[Install]
WantedBy=docker.service
